# cookbooks/appdynamics-machine-agent/templates/appd-machine.service.erb
[Unit]
Description=appd-machine-agent

[Service]
ExecStart=/opt/machineagent/bin/machine-agent

[Install]
WantedBy=multi-user.target
